# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Run pg_backuper every hour (smart scheduling in tool determines if backup is due)
# Output goes directly to stdout/stderr for Docker logs
0 * * * *	root	echo "=== Cron job triggered at $(date '+\%Y-\%m-\%d \%H:\%M:\%S') ===" && /usr/local/bin/pg_backuper $CONFIG_FILE 2>&1
